# systemd/teak-configurator-restart-service@.service
[Unit]
Description=Restart %i.service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart %i.service
RemainAfterExit=no
