# systemd/teak-configurator-restart-service@.path
# By enabling this unit for a service, we ensure that after the service starts systemd
# listens for touches on the configurator's restart file, and will run teak-configurator-restart-service@%i.service
# when the file is touched.
#
# This allows the configurator to restart systemd services without privileged access.
[Unit]
Description=Touch a file to restart %i.service

[Path]
PathModified=/run/teak-configurator/%i.restart

[Install]
WantedBy=%i.service
